# packages/lmaps/lmaps-1.0.0.tar.gz/lmaps-1.0.0/systemd/lmaps.service
[Unit]
Description=lMaPS Daemon
After=network.target

[Service]
PrivateTmp=true
Type=simple
Restart=always
RestartSec=10
ExecStart=/usr/bin/lmaps -D

[Install]
WantedBy=multi-user.target